# home/piupartsm/crontab
# m h dom mon dow (0|7=sun,1=mon)    command

# 
# notify about reboots
#
@reboot ddate | mail -s "piatti rebooted, please restart the piuparts-slave screen session" piupartss
#
# create reports once a day
#  (dinstall runs 1|7|13|19:52, so this is long after mirror pushes...)
#
5 0 * * * PYTHONPATH=/org/piuparts.debian.org/lib/python2.4/site-packages:/org/piuparts.debian.org/lib/python2.5/site-packages nice python /org/piuparts.debian.org/share/piuparts/piuparts-report  2>&1 | mail -s piuparts-report piupartsm

#
# monitor slave session every hour
#
23 * * * * /home/piupartsm/bin/detect_slave_problems 

#
# make sure slave is reminded daily (and also monitored if not started with ~piupartss/bin/slave_run)
#
42 0 * * * rm -f /org/piuparts.debian.org/master/monitor-slave/* ; cp /org/piuparts.debian.org/slave/screenlog.0 /org/piuparts.debian.org/master/monitor-slave/screenlog.new ; chmod -f 775 /org/piuparts.debian.org/master/monitor-slave/screenlog.new

#
# monitor for problems once a day
#
0  22 * * 0 /home/piupartsm/bin/detect_well_known_issues
45 22 * * * /home/piupartsm/bin/detect_stale_mounts
50 22 * * * /home/piupartsm/bin/detect_tmp_cruft
55 22 * * * /home/piupartsm/bin/detect_network_issues
0  23 * * * /home/piupartsm/bin/detect_archive_issues
05 23 * * * /home/piupartsm/bin/report_untestable_packages
10 23 * * * /home/piupartsm/bin/detect_well_known_errors

#
# misc
#
55 11 * * * /home/piupartsm/bin/prepare_backup

